# etc/systemd/system/postgresql-init.service
[Service]
Type=oneshot
ExecStart=/opt/bin/postgresql-init.sh
RemainAfterExit=true
